# mail-filter/openarc/files/openarc.initrc
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

CONFFILE="/etc/openarc/${RC_SVCNAME}.conf"
required_files="${CONFFILE}"

command="/usr/sbin/openarc"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-P ${pidfile} -x ${CONFFILE}"

depend() {
	use dns logger net
	before mta
}

