# mail-filter/openarc/files/openarc-initrc
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

extra_commands="checkconfig"
extra_started_commands="reload"

CONFFILE="/etc/openarc/${SVCNAME}.conf"
RUNDIR=/run/openarc
PIDFILE=$RUNDIR/openarc.pid

depend() {
	use dns logger net
	before mta
}

checkconfig() {
	/usr/sbin/openarc -c /etc/openarc/openarc.conf -n > /dev/null \
		|| return 1
}

start() {
	ebegin "Starting OpenARC"

	mkdir -m0750 -p $RUNDIR
	chown milter:milter $RUNDIR
	chmod g+s $RUNDIR

	start-stop-daemon --start --quiet \
		--exec /usr/sbin/openarc -- -c "${CONFFILE}" -P "${PIDFILE}"
	eend $?
}

stop() {
	ebegin "Stopping OpenARC"
	start-stop-daemon --stop --pidfile "${PIDFILE}"
	eend $?
}

reload() {
	checkconfig || return 1
	ebegin "Reloading ${SVCNAME}"
	start-stop-daemon --stop --signal HUP --quiet --pidfile $PIDFILE
	eend $?
}
